# service_template.service
[Unit]
Description=Script Description
After=network.target

[Service]
ExecStart=/path/to/script.sh

[Install]
WantedBy=default.target